# scripts/pxn-sub-stats.service
[Unit]
Description=PXN Shield subscription stats server
Documentation=https://github.com/NightRider322/PXN-SUB
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/pxn-sub/stats.env
ExecStart=/usr/bin/env python3 /usr/local/x-ui/pxn_stats.py
Restart=always
RestartSec=3
Nice=10
# light sandboxing — only needs to read /proc and write its output dir.
# install.sh adds a drop-in with ReadWritePaths for your theme directory.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/usr/local/x-ui

[Install]
WantedBy=multi-user.target
